# rootfs/etc/services.d/ttrss-updater/run
#!/usr/bin/with-contenv sh
# use with-contenv to access env vars

# run ttrss update daemon
s6-setuidgid www-data ${TTRSS_PHP_EXECUTABLE} /var/www/localhost/htdocs/update_daemon2.php
